# docker/conf/supervisor/services/horizon.conf.example
[program:laravel-horizon-$HORIZON_DOMAIN]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan horizon --domain=$HORIZON_DOMAIN
user=devilbox
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopwaitsecs=3600
logfile=/etc/supervisor/conf.d/horizon.log
